# tailscale/bin/tailscale-autostart.conf
# /etc/upstart/tailscale-autostart.conf
# Tailscale KUAL extension - autostart at boot
#
# This file is installed to /etc/upstart/ by the Enable Autostart menu item.
# The actual start decision is controlled by a trigger file on the userstore.

description "Tailscale autostart"

start on started volumd
stop on (stopping volumd or ota-update)

pre-start script
    BASEDIR="/mnt/us/extensions/tailscale"
    TRIGGER="${BASEDIR}/bin/autostart.enabled"
    SCRIPT="${BASEDIR}/bin/autostart_tailscale.sh"

    if [ -f "${TRIGGER}" ] && [ -f "${SCRIPT}" ]; then
        [ -x "${SCRIPT}" ] || chmod +x "${SCRIPT}"
        "${SCRIPT}" &
    fi

    exit 0
end script
